# tst/ubuntuDeploy/tst.service
[Unit]
Description=nodeApp
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/etc/init.d/start_tst.sh

[Install]
WantedBy=multi-user.target
